# conf.d/websockify.conf
[program:websockify]
command=websockify --web /novnc %(ENV_PUBPORT)s localhost:5900
autorestart=true
