# conf.d/websockify.conf
[program:websockify]
command=websockify --web /novnc 80 localhost:5900
autorestart=true
