# conf.d/xsession.conf
[program:xsession]
command=/etc/X11/Xsession
autorestart=true
