# conf.d/xvnc.conf
[program:xvnc]
command=Xvnc -ac -pn -dpms -depth 24 :0 -SecurityTypes=None -desktop=devrt
autorestart=true
