# conf.d/xvnc.conf
[program:xvnc]
command=Xvnc -ac -pn -depth 24 %(ENV_DISPLAY)s -rfbport=5900 -SecurityTypes=None -desktop=devrt
autorestart=true
